# workflows/client-intelligence/com.blueprint.ghl-conversations-pull.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.blueprint.ghl-conversations-pull</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/bp/projects/workflows/client-intelligence/ghl-conversations-pull.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>11</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>13</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>15</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
    </array>
    <key>StandardOutPath</key>
    <string>/Users/bp/projects/workflows/client-intelligence/launchd-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/bp/projects/workflows/client-intelligence/launchd-stderr.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/bp/projects/workflows/client-intelligence</string>
</dict>
</plist>
